# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7.4]
command = /usr/sbin/php-fpm7.4 -F
autostart = true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true